# scripts/backlight_capture@.service.in
[Unit]
Description=Backlight Capture Service
After=network.target

[Service]
ExecStart=/usr/bin/backlight_capture /dev/%I
StandardOutput=/var/log/backlight_capture.log
StandardError=/var/log/backlight_capture.log
Restart=on-failure
RestartSec=60s

[Install]
WantedBy=multi-user.target
