# scripts/com.backlight_capture.backlight_capture.plist.in
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.backlight_capture.backlight_capture</string>
    <key>ProgramArguments</key>
    <array>
      <string>/Library/Application Support/BacklightCapture/backlight_capture.sh</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <dict>
      <key>SuccessfulExit</key>
      <false/>
    </dict>
    <key>StandardOutPath</key>
    <string>/var/log/backlight_capture.log</string>
    <key>StandardErrorPath</key>
    <string>/var/log/backlight_capture.log</string>
</dict>
</plist>