# src/com.backlight_capture.backlight_capture.plist.in
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.backlight_capture.backlight_capture</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Library/Application Support/BacklightCapture/backlight_capture</string>
        <string>/dev/cu.usbserial-FTB6SPL3</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/tmp/backlight_capture.out</string>
    <key>StandardErrorPath</key>
    <string>/tmp/backlight_capture.err</string>
</dict>
</plist>